# crontab/shtc3.disable
@reboot python3 /home/aquapoly/monitoring/sensors/services/shtc3.py >> /var/log/crontab/cron.log 2>&1
* * * * * python3 /home/aquapoly/monitoring/sensors/services/shtc3.py >> /var/log/crontab/cron.log 2>&1

