# crontab/temp
@reboot cd /home/aquapoly/Aquadash/Sensors && python3 -m services.atlas_main 4 >> /var/log/crontab/cron.log 2>&1
* * * * * cd /home/aquapoly/Aquadash/Sensors && python3 -m services.atlas_main 4 >> /var/log/crontab/cron.log 2>&1
