# backend/deployment/backend-container.service
[Unit]
Description=Aquadash Backend Container
After=network.target camera-daemon.service
Requires=camera-daemon.service

[Service]
# Resolve host-dev GID dynamically so no hardcoded IDs are needed
ExecStart=/bin/bash -c '\
  GID=$(getent group host-dev | cut -d: -f3); \
  exec /usr/bin/docker run --rm \
    --name aquabackend \
    --user 1001:$GID \
    --group-add $GID \
    -v /run/camera:/run/camera:ro \
    backend-image:latest \
'
ExecStop=/usr/bin/docker stop aquabackend

KillMode=control-group
TimeoutStopSec=20
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
