# Sensors/crontab/ec
@reboot cd /home/aquadash/Documents/Aquadash/Sensors && python3 -m services.atlas_main 1 >> /var/log/cron.log 2>&1
* * * * * cd /home/aquadash/Documents/Aquadash/Sensors && python3 -m services.atlas_main 1 >> /var/log/cron.log 2>&1
