# Sensors/crontab/shtc3.disable
@reboot cd /opt/aquadash/sensors && python3 -m services.shtc3 >> /var/log/crontab/cron.log 2>&1
* * * * * cd /opt/aquadash/sensors && python3 -m services.shtc3 >> /var/log/crontab/cron.log 2>&1

