# Sensors/crontab/temp
@reboot cd /opt/aquadash/sensors && python3 -m services.atlas_main 3 >> /var/log/cron.log 2>&1
* * * * * cd /opt/aquadash/sensors && python3 -m services.atlas_main 3 >> /var/log/cron.log 2>&1
