# Sensors/crontab/ultrason.disable
@reboot cd /opt/aquadash/sensors && python3 -m services.ultrason_SEN0311 >> /var/log/crontab/cron.log 2>&1
* * * * * cd /opt/aquadash/sensors && python3 -m services.ultrason_SEN0311 >> /var/log/crontab/cron.log 2>&1
