# backend/deployment/devices/camera/camera-daemon.service
[Unit]
Description=Camera Daemon for /dev/video* devices
After=network.target
Wants=network.target

[Service]
# Ensure the runtime directory exists and has correct permissions
RuntimeDirectory=camera
RuntimeDirectoryMode=0750
User=root
Group=host-dev
ExecStart=/usr/bin/python3 __CAMERA_DIR__/camera_daemon.py
WorkingDirectory=__CAMERA_DIR__
Restart=on-failure
RestartSec=2s
# Ensure control socket is cleaned up on shutdown
ExecStopPost=/bin/rm -f /run/camera/control.sock

[Install]
WantedBy=multi-user.target
