# newBackend/src/features/sensors/material/crontab/camera.disable
@reboot python3 /home/aquapoly/monitoring/sensors/services/camera.py >> /var/log/crontab/cron.log 2>&1
* * * * * python3 /home/aquapoly/monitoring/sensors/services/camera.py >> /var/log/crontab/cron.log 2>&1
